# runit/record-weekly-shows/run
#!/usr/bin/env mksh
#
#
set -u -e -o pipefail

exec 2>&1
cd "$(realpath "$(dirname "$0")/../..")"
exec chpst -u da -U da scripts/record-weekly-shows

