# sv/nhk_high_2110/run
#!/bin/sh
#
#
/apps/dot_files/startups/livestreamer.sh high 2110
