# runit/recording-stream/run
#!/usr/bin/mksh
#
# This sets up the stream reserved for recording.
#

exec 2>&1
cd "$(realpath "$(dirname "$0")/../..")"
exec chpst -u da -U da scripts/stream high 2110
