# sv/nhk_record_weekly_shows/run
#!/usr/bin/mksh
#
#

set -u -e -o pipefail
PATH="$PATH:/apps/my_nhk/bin"
my_nhk record-next "HOKUSAI|SOMEWHERE|COOL JAPAN|UKIYO-E|Doki|Document|Buzz|Focus|Melo|Newsroom|Close-up|Inside"
