# runit/high/run
#!/usr/bin/mksh
#
#

exec 2>&1
cd "$(realpath "$(dirname "$0")/../..")"
exec chpst -u da -U da scripts/stream high 2112
